# client/model.conf
[program: <name>]
<calea catre executabil>

key_word1: <parametru1>
key_word2: <parametru2>
....