# client/test1.conf
[program:cat]
/bin/cat
stdout: out.txt create append
stderr: 
stdin:
owner:
group_owner:
umask:
priority:
directory: