# configure/test_sleep.conf
[program: sleep]
/home/claudia/Documents/Supervisor/client/sleep
stdout: out.txt create append
directory: /home/claudia/Documents
owner: claudia
group_owner: claudia
priority: 4
