# configure/test_ls.conf
[program: ls]
/bin/ls
directory: /home/claudia/Documents
owner: claudia
group_owner: claudia
priority: 4
