# configure/test1.conf
[program: test]
/home/claudia/Documents/Supervisor/client/test1
stdout: out.txt create append
stdin: in.txt
owner: claudia
group_owner: claudia
umask: 0024
priority: 4
