# configure/test2.conf
[program: test2]
/home/claudia/Documents/Supervisor/client/test1
stdin: in.txt
owner: claudia
group_owner: claudia
umask: 0024
priority: 4
